# install/server/adminset.service
[Unit]
Description=Adminset Web Service
After=network.target
Wants=redis.service mariadb.service mongodb.service

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/var/opt/adminset/main
Environment=PYTHONPATH=/var/opt/adminset/main
ExecStart=/usr/bin/env python3 /var/opt/adminset/main/manage.py runserver 0.0.0.0:8000
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target